# releng/centos-6-i386-build/etc/rc.local
../../centos-6-x86_64-build/etc/rc.local